# docker/stack.conf
[program:monkey]
command=/var/monkey/monkey_island/ubuntu/systemd/start_server.sh
